# nextcloud18/systemd/nextcloud18@.timer
[Unit]
Description=Run Nextcloud cron.php every 15 minutes (%i)

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Unit=nextcloud18@%i.service

[Install]
WantedBy=timers.target
